# systemd/codex-pm-service.service
# codex-pm-service — Bot job management service (PM-D5, PM-D6)
#
# Lightweight per-user service that manages bot runs for work items.
# Listens on $XDG_RUNTIME_DIR/codex-pm.sock via JSON-RPC-lite (PM-D7, PM-D8).
#
# Locked decisions:
# - D136: Unix-domain-socket IPC (socket activation supported)
# - D141: start on login (`WantedBy=default.target`)
# - D142: stay running while logged in (idle is acceptable; no always-processing loops)
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp codex-pm-service.{service,socket} ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now codex-pm-service.service
#
# Usage:
#   systemctl --user status codex-pm-service
#   journalctl --user -u codex-pm-service -f

[Unit]
Description=Codex PM Bot Job Management Service
Documentation=https://github.com/anthropics/codex/tree/main/docs/SPEC-PM-003-bot-system
Requires=codex-pm-service.socket
After=codex-pm-service.socket

[Service]
Type=simple
# Receive the listening socket from codex-pm-service.socket (D136).
Sockets=codex-pm-service.socket
ExecStart=%h/.cargo/bin/codex-pm-service
Restart=on-failure
RestartSec=5

# Environment
Environment=RUST_LOG=info
Environment=XDG_RUNTIME_DIR=%t

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t %h/.local/share/codex-pm

# Resource limits (lightweight service)
MemoryMax=256M
TasksMax=64

[Install]
WantedBy=default.target
