# systemd/codex-pm-service.service
# codex-pm-service — Bot job management service (PM-D5, PM-D6)
#
# Lightweight per-user service that manages bot runs for work items.
# Listens on $XDG_RUNTIME_DIR/codex-pm.sock via JSON-RPC-lite (PM-D7, PM-D8).
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp codex-pm-service.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable codex-pm-service
#
# Usage:
#   systemctl --user start codex-pm-service
#   systemctl --user status codex-pm-service
#   journalctl --user -u codex-pm-service -f
#
# The service is intended to be **exit-when-idle** (D135). In v1, prefer
# systemd socket activation (D136) so the service only runs while it has work.
#
# NOTE: This unit is a development baseline; production hardening (socket
# activation + auto-resume semantics) is tracked in SPEC-PM-003 Phase-1.

[Unit]
Description=Codex PM Bot Job Management Service
Documentation=https://github.com/anthropics/codex/tree/main/docs/SPEC-PM-003-bot-system
After=default.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/codex-pm-service
Restart=on-failure
RestartSec=5

# Environment
Environment=RUST_LOG=info
Environment=XDG_RUNTIME_DIR=%t

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t

# Resource limits (lightweight service)
MemoryMax=256M
TasksMax=64

[Install]
WantedBy=default.target
