# systemd/codex-pm-resume.timer
# codex-pm-resume — Auto-resume incomplete bot runs after reboot/login
#
# Fires once, 30 seconds after the user session starts, to trigger
# socket-activation of codex-pm-service.  The service's startup path
# calls resume_incomplete() before accepting connections, so any
# incomplete runs left on disk are automatically re-executed.
#
# Combined with exit-when-idle (D135), the service shuts down after
# 60 seconds if no real client connects.
#
# Install (one-time):
#   mkdir -p ~/.config/systemd/user
#   cp codex-pm-resume.{timer,service} ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now codex-pm-resume.timer

[Unit]
Description=Codex PM Auto-Resume Timer
After=codex-pm-service.socket

[Timer]
OnStartupSec=30s

[Install]
WantedBy=timers.target
