# systemd/codex-pm-resume.service
# codex-pm-resume.service — Oneshot that pokes the PM service socket
#
# Triggered by codex-pm-resume.timer after login.  Connects to the
# socket (triggering socket-activation of codex-pm-service), sends a
# hello handshake, then exits.  The main service resumes incomplete
# runs at startup and idles out via D135 if no real client connects.
#
# This unit has no [Install] section — it is started only by the timer.

[Unit]
Description=Codex PM Auto-Resume Trigger
After=codex-pm-service.socket
Requires=codex-pm-service.socket

[Service]
Type=oneshot
ExecStart=%h/.cargo/bin/codex-pm-service --ping
TimeoutStartSec=30
Environment=XDG_RUNTIME_DIR=%t
